# recipes-application/nltsbeac/nltsbeac/nltsbeac.service
[Unit]
Description=Main nlts application

[Service]
ExecStart = /usr/bin/nltsbeac
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

