# recipes-application/ldms/ldms-2.4.0/ldms.service
[Unit]
Description=Main ldms application

[Service]
ExecStart = /usr/bin/ldms 
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
