# recipes-application/ldms-cfg/ldms-cfg-1.1/ldms-cfg.service
[Unit]
Description=Setting up  GPIO sysfs interface

[Service]
ExecStart = /usr/bin/ldms-cfg

[Install]
WantedBy=multi-user.target
